# .config/autostart/start_apps.desktop
[Desktop Entry]
Exec=/home/viktor_mitin/configs/start_apps.sh
Name=start_apps
Type=Application
Version=1.0
