# .config/autostart/start_apps.desktop
[Desktop Entry]
Exec=/home/c/configs/start_apps.sh
Name=start_apps
Type=Application
Version=1.0
Comment=
RunHook=0

